# worker/cron/crontab
*/2 * * * * python3 task_runner.py
*/2 * * * * echo "Hello world" >> /var/log/cron.log 2>&1

